# viewer/systemd/s9-knowledge-viewer.service
[Unit]
Description=S9 Knowledge Viewer
After=network.target
# La release activa se resuelve SIEMPRE a través del symlink `current`; esta unit
# no referencia ningún layout legacy.
AssertPathExists=/opt/s9-knowledge/current/viewer
AssertPathExists=/etc/s9-knowledge/viewer.env

[Service]
Type=simple
WorkingDirectory=/opt/s9-knowledge/current/viewer
EnvironmentFile=/etc/s9-knowledge/viewer.env
ExecStart=/opt/s9-knowledge/current/viewer/.venv/bin/uvicorn app.main:app --host ${S9K_VIEWER_HOST} --port ${S9K_VIEWER_PORT}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
