# viewer/systemd/s9-knowledge-healthcheck.timer
[Unit]
Description=Run S9 Knowledge operational healthcheck every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
