# viewer/systemd/s9-knowledge-healthcheck.timer
[Unit]
Description=Run S9 Knowledge operational healthcheck every 5 minutes
# El servicio es Type=oneshot y OnUnitActiveSec cuenta desde la ultima
# activacion: mientras una ejecucion siga viva no se dispara otra, asi que los
# checks no se solapan.

[Timer]
Unit=s9-knowledge-healthcheck.service
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
# Persistent solo tiene efecto con OnCalendar; aqui no aporta nada.

[Install]
WantedBy=timers.target
